# sinabung-monitoring.service
[Unit]
Description=Sinabung Monitoring & DevOps Dashboard
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/var/www/asetpedia/sinabung-monitoring
EnvironmentFile=-/var/www/asetpedia/sinabung-monitoring/.env
Environment="PATH=/var/www/asetpedia/sinabung-monitoring/venv/bin"
Environment="DEV_MODE=False"

# Run Gunicorn: 1 worker with 4 threads for high-density API monitoring
ExecStart=/var/www/asetpedia/sinabung-monitoring/venv/bin/gunicorn --workers 1 --threads 4 --bind 0.0.0.0:9000 app:app

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
